# stage3_dronemap/11-readonly/files/rc.local
/usr/sbin/ntpdate -b time.google.com # change the ntp server according to your location
